# service/mistserver/run
#!/bin/bash
echo "n" | MistController -c /config/server.conf
